# teamcity-agent-code/roles/teamcityagent/templates/teamcity-agent.conf
[program:{{application}}]
command=/{{username}}/{{application}}/data/bin/agent.sh start
environment=HOME="/home/{{username}}",USER="{{username}}"
user={{username}}
autostart=true
